# Sidora-Batch/config/sidora-batch.service
[Unit]
Description=Sidora Batch Service
After=syslog.target ssd.mount networking.service

[Service]
User=fedora
ExecStart=@project.build.directory@/@project.build.finalName@.jar
SuccessExitStatus=143
SyslogIdentifier=sidora-batch-service
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target